# packages/sollol/sollol-0.9.48.tar.gz/sollol-0.9.48/systemd/sollol-gpu-reporter.service
[Unit]
Description=SOLLOL GPU Reporter Service
After=network.target redis.service
Documentation=https://github.com/BenevolentJoker-JohnL/SynapticLlamas
Wants=redis.service

[Service]
Type=simple
User=%u
WorkingDirectory=%h/SOLLOL
ExecStart=/usr/bin/env python3 %h/SOLLOL/gpu_reporter.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables (override in /etc/systemd/system/sollol-gpu-reporter.service.d/override.conf)
Environment="REDIS_HOST=localhost"
Environment="REDIS_PORT=6379"
Environment="OLLAMA_PORT=11434"
Environment="REPORT_INTERVAL=5"

[Install]
WantedBy=default.target
